# tools/systemd_user/vivid-hyprland.target.in
[Unit]
Description=Vivid Hyprland runtime
Wants=vivid-producer.service vivid-webui.service vivid-wayland-hyprland.service
After=graphical-session.target
